# lib/capper/templates/delayed_job@.service.erb
[Unit]
Description=Delayed Job Worker (%I)

[Service]
Type=forking
PIDFile=<%= pid_path %>/delayed_job.%I.pid
ExecStart=<%= delayed_job_script %> start %I
ExecStop=<%= delayed_job_script %> stop %I
Restart=always

[Install]
WantedBy=default.target
