# lib/capper/templates/thin.service.erb
[Unit]
Description=Thin Server

[Service]
Type=forking
PIDFile=<%= thin_pidfile %>
ExecStart=<%= thin_script %> start
ExecStop=/bin/kill -QUIT $MAINPID

[Install]
WantedBy=default.target
